# deploy/bundle-supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
pidfile=/tmp/supervisord.pid

[program:postgres]
command=/usr/local/bin/docker-entrypoint.sh postgres
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

[program:redis]
command=redis-server --bind 127.0.0.1 --port %(ENV_REDIS_PORT)s --save 60 1 --appendonly yes --appendfsync everysec
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20

[program:sub2api]
command=sh -c 'until pg_isready -U "${POSTGRES_USER:-sub2api}" -d "${POSTGRES_DB:-sub2api}" >/dev/null 2>&1; do sleep 1; done; until redis-cli -p "${REDIS_PORT:-6379}" ping 2>/dev/null | grep -q PONG; do sleep 1; done; exec /app/docker-entrypoint.sh /app/sub2api'
user=sub2api
directory=/app
autostart=true
autorestart=true
startsecs=15
startretries=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=30
